# deploy/spread-scanner-autodeploy.timer
[Unit]
Description=Check GitHub for Spread Scanner updates

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=15s
Persistent=true
Unit=spread-scanner-autodeploy.service

[Install]
WantedBy=timers.target
